# 2026/installers/debian/ao-exchange.service
[Unit]
Description=AO Exchange - Assign Onward exchange agent
After=network.target

[Service]
Type=simple
User=ao-exchange
Group=ao-exchange
WorkingDirectory=/var/lib/ao-exchange
ExecStart=/usr/bin/ao-exchange run /etc/ao-exchange/exchange.toml
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/ao-exchange
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ao-exchange

[Install]
WantedBy=multi-user.target
